# files/etc/init/senioritis.conf
description "Senioritis"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /home/ngoke/Code/senioritis/start.sh
